# server/systemd/srt-relay-watchdog.timer
[Unit]
Description=Run SRT Relay Watchdog every 60s

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
