# server/systemd/srtla-ip-mapper.timer
[Unit]
Description=SRTLA IP Mapper timer (refresh every 10s)

[Timer]
OnBootSec=10
OnUnitActiveSec=10

[Install]
WantedBy=timers.target
